# config/supervisord.conf
; ==============================================================================
; DebianXRDP-Pro
; Supervisor Configuration
; ==============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid

; ------------------------------------------------------------------------------
; DBus
; ------------------------------------------------------------------------------

[program:dbus]

command=/usr/bin/dbus-daemon --system --nofork

priority=10

autostart=true

autorestart=true

startsecs=3

stdout_logfile=/var/log/dbus.log

stderr_logfile=/var/log/dbus-error.log

; ------------------------------------------------------------------------------
; XRDP Session Manager
; ------------------------------------------------------------------------------

[program:xrdp-sesman]

command=/usr/sbin/xrdp-sesman --nodaemon

priority=20

autostart=true

autorestart=true

startsecs=3

stdout_logfile=/var/log/xrdp-sesman.log

stderr_logfile=/var/log/xrdp-sesman-error.log

; ------------------------------------------------------------------------------
; XRDP
; ------------------------------------------------------------------------------

[program:xrdp]

command=/usr/sbin/xrdp --nodaemon

priority=30

autostart=true

autorestart=true

startsecs=3

stdout_logfile=/var/log/xrdp.log

stderr_logfile=/var/log/xrdp-error.log

; ------------------------------------------------------------------------------
; Event Listener
; ------------------------------------------------------------------------------

[eventlistener:process-monitor]

command=/bin/cat

events=PROCESS_STATE

buffer_size=100
